# deploy/augmentorai.service
# AugmentorAI — FastAPI backend service
# Install: sudo cp deploy/augmentorai.service /etc/systemd/system/
# Then:    sudo systemctl daemon-reload
#          sudo systemctl enable --now augmentorai
#
# Runs as opc user, using the project's own venv.

[Unit]
Description=AugmentorAI — AI-powered interview practice copilot API
After=network-online.target ollama.service
Wants=network-online.target ollama.service

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/home/opc/projects/AugmentorAI/server

# Activate venv and start uvicorn
ExecStart=/home/opc/projects/AugmentorAI/server/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8010 --workers 2

Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/opc/projects/AugmentorAI/server/.env

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/opc/projects/AugmentorAI/server/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
